# packages/agent/systemd/usipipo-agent.service
[Unit]
Description=uSipipo VPN Agent
After=network.target

[Service]
Type=simple
User=usipipo
Group=usipipo
WorkingDirectory=/opt/usipipo-agent
EnvironmentFile=/opt/usipipo-agent/.env
ExecStart=/opt/usipipo-agent/usipipo-agent
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=usipipo-agent

# Allow WireGuard operations via netlink
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW

[Install]
WantedBy=multi-user.target
