# apps/backend/deploy/usipipo-backend.service
# =============================================================================
# uSipipo Backend - Systemd Service
# =============================================================================
# Production-ready service configuration with security hardening
#
# BEFORE INSTALLING:
# 1. Copy this file to /etc/systemd/system/usipipo-backend.service
# 2. Replace all placeholders (marked with <PLACEHOLDER>) with your values
# 3. Run: sudo systemctl daemon-reload
# 4. Run: sudo systemctl enable --now usipipo-backend
# =============================================================================

[Unit]
Description=uSipipo Backend API Service
Documentation=https://github.com/uSipipo-Team/usipipo-backend
After=network.target postgresql.service redis.service
Wants=postgresql.service redis.service

[Service]
Type=exec

# User and Group
# Replace with your system user/group (use 'id <username>' to verify)
User=mowgli
Group=mowgli

# Working Directory
# Update to your actual deployment path
WorkingDirectory=/home/mowgli/usipipo/usipipo-backend

# Environment Configuration
Environment="PATH=/home/mowgli/usipipo/usipipo-backend/.venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="VIRTUAL_ENV=/home/mowgli/usipipo/usipipo-backend/.venv"
EnvironmentFile=/home/mowgli/usipipo/usipipo-backend/.env

# Main Process
# FastAPI backend with uvicorn workers
# Adjust --port and --workers based on your deployment needs
ExecStart=/home/mowgli/usipipo/usipipo-backend/.venv/bin/python -m uvicorn src.main:app \
    --host 0.0.0.0 \
    --port 8001 \
    --workers 2 \
    --access-log

# Graceful Shutdown
# Allows in-flight requests to complete before termination
TimeoutStopSec=30
KillMode=mixed
KillSignal=SIGTERM

# Restart Policy
# Always restart on failure or clean exit
Restart=always
RestartSec=10

# Security Hardening
# See systemd.exec(5) for detailed documentation
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=false
ReadWritePaths=/home/mowgli/usipipo/usipipo-backend/logs

# Resource Limits
# Prevent resource exhaustion attacks
LimitNOFILE=65536
LimitNPROC=4096

# Logging
# Output goes to systemd journal (view with: journalctl -u usipipo-backend)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=usipipo-backend

[Install]
WantedBy=multi-user.target
